# java/demo/service.supervisord.conf
[program:demo]
command=/opt/apache-tomcat-7.0.67/bin/supervisor-tomcat-wrapper.sh